# 2024-01-XX_GrafanaLab_demo/demo/client/supervisor.d/app.ini
[program:app]
directory=/app
command=java
    -Dotel.service.name=%(ENV_SERVICE_NAME)s
    -Dotel.resource.attributes.service.name=%(ENV_SERVICE_NAME)s
    -Dotel.resource.attributes.application=%(ENV_SERVICE_NAME)s
    -Dotel.resource.attributes.host=%(host_node_name)s
    -Dotel.resource.attributes.env=%(ENV_SERVICE_ENV)s
    -Dotel.instrumentation.micrometer.enabled=false
    -jar /app/main.jar --spring.application.name=%(ENV_SERVICE_NAME)s --management.otlp.metrics.export.resourceAttributes.service.name=%(ENV_SERVICE_NAME)s --management.otlp.metrics.export.resourceAttributes.service.instance.id=%(host_node_name)s --management.metrics.tags.env=%(ENV_SERVICE_ENV)s --management.metrics.tags.application=%(ENV_SERVICE_NAME)s --debug
autorestart=false
startretries=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0