# 2023-11-23_Vector_hands-on/demo/app/supervisor/supervisor.d/vector.ini
[program:vector]
directory=/workspace/vector
command=/workspace/vector/vector --log-format=json --config /workspace/vector/config/**/*.toml
autorestart=false
startretries=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0